# ubuntu/90_myrule.rules
ACTION=="add", ATTRS{name}=="Kensington Expert Wireless TB", RUN+="/usr/local/bin/remap.sh"

